# gunicorn.conf
[program:gunicorn]

command=/home/python/projects/students_tracker/env/bin/gunicorn students_tracker.wsgi
directory=/home/python/projects/students_tracker/src
autostart=true
autorestart=true
